# package/carbonio-message-dispatcher-auth-sidecar.service
[Unit]
Description=Carbonio Message Dispatcher auth sidecar proxy
Documentation=https://docs.zextras.com/
Requires=network-online.target
After=network-online.target
PartOf=carbonio.target

[Service]
User=carbonio-message-dispatcher
ExecStart=/usr/bin/consul connect envoy \
    -token-file /etc/carbonio/message-dispatcher/service-discover/token \
    -admin-bind localhost:0 \
    -sidecar-for carbonio-message-dispatcher-auth
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=15
KillSignal=SIGINT
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
